# systemd-user/offlineimap-notify.service
[Unit]
Description=Start offlineimap as a daemon

[Service]
ExecStart=/usr/bin/offlineimap-notify

[Install]
WantedBy=default.target
